# cron_40000631-master/src/crontab
*/5 * * * * curl -s http://periodictest.40000631.qpc.hal.davecutting.uk/ > /dev/null
